# crontab/check-ping-state.sh
#!/bin/bash
#
# Script para ejecutar aplicacion Node check-ping-state
#

# 1. Carga el PATH donde se encuentra check-ping-state
cd /path/to/node/check-ping-state/
# 2. Ejecuta la aplicacion de Node check-ping-state
node app.js
