# crontab/crontab.conf
#
# Configuracion en crontab
# 
# Para facilitar la configuracion en crontab:
# https://crontab-generator.org/
#


#
# Ejecuta tarea para comprobar el estado de respuesta de una peticion PING
#
# Cada 5 minutos
#
*/5 * * * * /path/to/script/check-ping-state.sh
